# deploy/vigil-worker.service
[Unit]
Description=Vigil remote scan worker
After=network-online.target
[Service]
User=__USER__
WorkingDirectory=__HOME__/Vigil
Environment=KRONOS_REPO_PATH=__HOME__/Kronos
Environment=KRONOS_DEVICE=cpu
Environment=VIGIL_WORKER_PORT=8090
ExecStart=__HOME__/Vigil/.venv/bin/python -m scanner.worker
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
