# deploy/vigil-portfolio.timer
[Unit]
Description=Vigil portfolio guard schedule (free-time checks)
[Timer]
OnCalendar=Mon-Fri 12:00 Europe/London
OnCalendar=Mon-Fri 16:00 Europe/London
Persistent=true
[Install]
WantedBy=timers.target
