# deploy/vigil-signals.timer
[Unit]
Description=Vigil pre-market sweep (Europe first, then global/US baskets)
[Timer]
# One sweep a day, well before the London open. London tz tracks BST/GMT.
# Order of markets is set by SIGNAL_MARKETS (recommend: europe,global-liquid).
OnCalendar=Mon-Fri 05:00 Europe/London
Persistent=true
[Install]
WantedBy=timers.target
