# deploy/vigil-signals.timer
[Unit]
Description=Run Vigil signals on a schedule (weekday pre-open, UTC)
[Timer]
# 11:00 UTC weekdays = before the US open; add more OnCalendar lines for more runs.
OnCalendar=Mon-Fri 11:00
Persistent=true
RandomizedDelaySec=300
[Install]
WantedBy=timers.target
